# deploy/tspl-simulator.service
[Unit]
Description=TSPL Simulator Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/tspl-simulator/backend
EnvironmentFile=-/opt/tspl-simulator/backend/.env
ExecStartPre=/bin/sh -c 'mkdir -p /opt/tspl-simulator/backend/data/API_print /opt/tspl-simulator/backend/data/MQTT_print'
ExecStart=/opt/tspl-simulator/backend/tspl-simulator
Restart=on-failure
RestartSec=5s

# 安全設置
NoNewPrivileges=true
PrivateTmp=true

# 日誌設置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tspl-simulator

[Install]
WantedBy=multi-user.target
